# hourly-wait-reporter.service.example
[Unit]
Description=check_price hourly wait-level reporter
After=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/shu2891/check_price
EnvironmentFile=/home/shu2891/.config/check_price.env
ExecStart=/usr/bin/python3 /home/shu2891/check_price/hourly_wait_reporter.py --config /home/shu2891/check_price/watchlist.json --loop-minutes 240 --llama off
Restart=always
RestartSec=10

[Install]
WantedBy=default.target
